# tools/touchscreen/hx83121a-touch-recovery.service
[Unit]
Description=HX83121A Touchscreen Recovery (Boot ROM trigger)
After=sys-bus-i2c-devices-4\x2d004f.device
Before=graphical.target
ConditionPathExists=/sys/bus/i2c/devices/4-004f

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hx83121a-touch-recovery
RemainAfterExit=yes
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
